# k3s_for_yc/vector/vector.service
[Service]
ExecStart=/usr/bin/vector --config /var/lib/vector/vector.yaml
Restart=always
RestartSec=5
Environment=PATH=/usr/bin:/bin:/usr/local/bin
